# webhook.conf
[program:webhook]
command=/usr/local/bin/webhook
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/webhook.err.log
stdout_logfile=/var/log/supervisor/webhook.out.log